# docker/php/etc/supervisor/conf.d/log.conf
[program:symfony-dev.log]
command=/usr/bin/tail -f /app/var/log/dev.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true