# docker/php/etc/supervisor/conf.d/php-fpm.conf
[program:php]
 command=docker-php-entrypoint php-fpm # This is our main process so let us put it here
